# debian/deb-src/etc/cron.d/geekbot
*/10 * * * * root /usr/bin/sync-git /var/lib/geekbot
